# rc.local
#!/bin/sh

apt update >/dev/null 2>&1 &

#Backup all the configuration file with noninteractive shell BUG
#/usr/local/bin/bak_cfg.sh backup &  

