# rc.local
#!/bin/sh

apt update >/dev/null 2>&1 &

#BUG# Backup all the configuration file with noninteractive shell
#/usr/local/bin/bak_cfg.sh backup &  
